# base/root_part/usr/local/bin/expand.sh
#!/bin/sh
#expand / partition to full sd card

PART_NUM='2'

# Get the starting offset of the root partition
START_SECTOR=$(parted /dev/mmcblk0 -ms unit s p | grep "^${PART_NUM}" | cut -f 2 -d: | sed 's/[^0-9]//g')

#repartition SD Card
sudo fdisk /dev/mmcblk0 >> /dev/null 2>&1 <<EOF
d
$PART_NUM
p
n
p
$PART_NUM
$START_SECTOR

p
w
q
EOF

#expand filesystem on full sd card

# now set up an init.d script
cat <<\EOF > /etc/init.d/resize2fs_once &&
#!/bin/sh
### BEGIN INIT INFO
# Provides:          resize2fs_once
# Required-Start:
# Required-Stop:
# Default-Start: 3
# Default-Stop:
# Short-Description: Resize the root filesystem to fill partition
# Description:
### END INIT INFO

. /lib/lsb/init-functions

case "$1" in
  start)
    log_daemon_msg "Starting resize2fs_once" &&
    resize2fs /dev/mmcblk0p2 &&
    update-rc.d resize2fs_once remove &&
    rm /etc/init.d/resize2fs_once &&
    log_end_msg $?
    ;;
  *)
    echo "Usage: $0 start" >&2
    exit 3
    ;;
esac
EOF
  chmod +x /etc/init.d/resize2fs_once &&
  update-rc.d resize2fs_once defaults


echo -n 'TRUE'
#done!
